# udev.rules
# EC200G Modem - Load option driver
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="2c7c", ATTRS{idProduct}=="0904", \
RUN+="/sbin/modprobe option", \
RUN+="/sbin/modprobe cdc_ether", \
RUN+="/bin/sh -c 'echo 2c7c 0904 > /sys/bus/usb-serial/drivers/option1/new_id'"

# EC200G - Bind network interfaces to cdc_ether  
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="2c7c", ATTRS{idProduct}=="0904", \
ATTR{bInterfaceNumber}=="00", \
RUN+="/bin/sh -c 'sleep 2; \
DEVPATH=$(echo %p | grep -oE \"[0-9]-[0-9]+\\.[0-9]+\"); \
echo $DEVPATH:1.0 > /sys/bus/usb/drivers/option/unbind 2>/dev/null || true; \
echo $DEVPATH:1.1 > /sys/bus/usb/drivers/option/unbind 2>/dev/null || true; \
echo $DEVPATH:1.0 > /sys/bus/usb/drivers/cdc_ether/bind 2>/dev/null || true; \
echo $DEVPATH:1.1 > /sys/bus/usb/drivers/cdc_ether/bind 2>/dev/null || true'"

# EC200G Serial ports
KERNEL=="ttyUSB*", ATTRS{idVendor}=="2c7c", ATTRS{idProduct}=="0904", ATTRS{bInterfaceNumber}=="02", SYMLINK+="modem_dm"
KERNEL=="ttyUSB*", ATTRS{idVendor}=="2c7c", ATTRS{idProduct}=="0904", ATTRS{bInterfaceNumber}=="03", SYMLINK+="modem_at"
KERNEL=="ttyUSB*", ATTRS{idVendor}=="2c7c", ATTRS{idProduct}=="0904", ATTRS{bInterfaceNumber}=="04", SYMLINK+="modem_ppp"
KERNEL=="ttyUSB*", ATTRS{idVendor}=="2c7c", ATTRS{idProduct}=="0904", ATTRS{bInterfaceNumber}=="05", SYMLINK+="modem_gps"

# RS485
KERNEL=="ttyUSB*", ATTRS{idVendor}=="0403", ATTRS{idProduct}=="6001", SYMLINK+="rs485"